# firmware/yocto/meta-edge-secure/recipes-core/bootstrap/edge-bootstrap/edge-bootstrap.service
[Unit]
Description=Edge AI First Boot Bootstrap
After=network-online.target docker.service data.mount time-sync.target
Wants=network-online.target time-sync.target
Requires=docker.service data.mount

[Service]
Type=oneshot
ExecStart=/usr/bin/edge-bootstrap.sh
RemainAfterExit=yes
StandardOutput=journal+console
StandardError=journal+console
TimeoutStartSec=600

[Install]
WantedBy=multi-user.target

